# systemfiles.archLinux/etc/cron.daily/fkol-k4.backup
#!/bin/bash

#   File: /et/cron.daily/fkol-k4.backup
#   Type: Bash script

#   Backup user files (daily) to '/Backup/' disk
#   Copyright (C) 2016 fkol-k4 (fkolyvas at gmail dot com)

#   This program is free software; you can redistribute it and/or modify
#   it under the terms of the GNU General Public License as published by
#   the Free Software Foundation version 2.

#   This program is distributed in the hope that it will be useful,
#   but WITHOUT ANY WARRANTY; without even the implied warranty of
#   MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
#   GNU General Public License for more details.

#   You should have received a copy of the GNU General Public License along
#   with this program; if not, write to the Free Software Foundation, Inc.,
#   51 Franklin Street, Fifth Floor, Boston, MA 02110-1301 USA.


function backup_music() {
    src="/Files-WD/Music"
    dst="/mnt/BackupFiles/Music"
    logfile="$src/rsync-backup.log"
    if [[ -d $src ]] && [[ -d $dst ]]; then
        rsync -aqz --log-file=$logfile $src/ $dst/
        sync
    fi
}

function backup_keys() {
    src="/Files-WD/keys"
    dst="/mnt/BackupFiles/keys"
    logfile="$src/rsync-backup.log"
    if [[ -d $src ]] && [[ -d $dst ]]; then
        rsync -aqz --log-file=$logfile $src/ $dst/
        sync
    fi
}

function backup_documents() {
    src="/Files-WD/Documents"
    dst="/mnt/BackupFiles/Documents"
    logfile="$src/rsync-backup.log"
    if [[ -d $src ]] && [[ -d $dst ]]; then
        rsync -aqz --log-file=$logfile $src/ $dst/
        sync
    fi
}

function main() {
    # main function
    #message=""
    #clear; echo "$message"; sleep 4
    backup_keys
    backup_documents
    backup_music
    sync
}

#main "$@"
